# apps/playwright-mcp/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:playwright-mcp]
command=/usr/local/bin/neko-playwright-mcp --port=%(ENV_PLAYWRIGHT_MCP_PORT)s --host=0.0.0.0 --cdp-endpoint=%(ENV_PLAYWRIGHT_MCP_CDP_ENDPOINT)s
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:cliproxy]
command=/usr/local/bin/neko-cliproxy
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
