# apps/cloakbrowser/supervisord.conf
[program:cloakbrowser-proxy-auth-forwarder]
environment=CLOAKBROWSER_PROXY_SERVER="%(ENV_CLOAKBROWSER_PROXY_SERVER)s",CLOAKBROWSER_PROXY_USERNAME="%(ENV_CLOAKBROWSER_PROXY_USERNAME)s",CLOAKBROWSER_PROXY_PASSWORD="%(ENV_CLOAKBROWSER_PROXY_PASSWORD)s",CLOAKBROWSER_LOCAL_PROXY_HOST="%(ENV_CLOAKBROWSER_LOCAL_PROXY_HOST)s",CLOAKBROWSER_LOCAL_PROXY_PORT="%(ENV_CLOAKBROWSER_LOCAL_PROXY_PORT)s"
command=/usr/local/bin/python3 /usr/local/bin/neko-cloakbrowser-proxy-auth-forwarder.py
autorestart=true
priority=790
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cloakbrowser-proxy-auth-forwarder.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:cloakbrowser]
environment=TZ="%(ENV_TZ)s",HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",CLOAKBROWSER_PROFILE_DIR="%(ENV_CLOAKBROWSER_PROFILE_DIR)s",CLOAKBROWSER_CDP_PORT="%(ENV_CLOAKBROWSER_CDP_PORT)s",CLOAKBROWSER_PROXY_PORT="%(ENV_CLOAKBROWSER_PROXY_PORT)s",CLOAKBROWSER_WINDOW_SIZE="%(ENV_CLOAKBROWSER_WINDOW_SIZE)s",CLOAKBROWSER_REMOTE_DEBUGGING_ADDRESS="%(ENV_CLOAKBROWSER_REMOTE_DEBUGGING_ADDRESS)s",CLOAKBROWSER_EXTRA_ARGS="%(ENV_CLOAKBROWSER_EXTRA_ARGS)s",CLOAKBROWSER_START_URL="%(ENV_CLOAKBROWSER_START_URL)s",CLOAKBROWSER_BIN="%(ENV_CLOAKBROWSER_BIN)s",CLOAKBROWSER_STEALTH_ARGS_ENABLED="%(ENV_CLOAKBROWSER_STEALTH_ARGS_ENABLED)s",CLOAKBROWSER_FINGERPRINT_PLATFORM="%(ENV_CLOAKBROWSER_FINGERPRINT_PLATFORM)s",CLOAKBROWSER_FINGERPRINT="%(ENV_CLOAKBROWSER_FINGERPRINT)s",CLOAKBROWSER_DESKTOP_POINTER_ARGS_ENABLED="%(ENV_CLOAKBROWSER_DESKTOP_POINTER_ARGS_ENABLED)s",CLOAKBROWSER_DESKTOP_GUARD_ENABLED="%(ENV_CLOAKBROWSER_DESKTOP_GUARD_ENABLED)s",CLOAKBROWSER_DESKTOP_GUARD_SCRIPT="%(ENV_CLOAKBROWSER_DESKTOP_GUARD_SCRIPT)s",CLOAKBROWSER_PROXY_SERVER="%(ENV_CLOAKBROWSER_PROXY_SERVER)s",CLOAKBROWSER_PROXY_USERNAME="%(ENV_CLOAKBROWSER_PROXY_USERNAME)s",CLOAKBROWSER_PROXY_PASSWORD="%(ENV_CLOAKBROWSER_PROXY_PASSWORD)s",CLOAKBROWSER_PROXY_BYPASS_LIST="%(ENV_CLOAKBROWSER_PROXY_BYPASS_LIST)s",CLOAKBROWSER_LOCAL_PROXY_HOST="%(ENV_CLOAKBROWSER_LOCAL_PROXY_HOST)s",CLOAKBROWSER_LOCAL_PROXY_PORT="%(ENV_CLOAKBROWSER_LOCAL_PROXY_PORT)s"
command=/usr/local/bin/neko-cloakbrowser-entrypoint.sh
stopsignal=INT
autorestart=true
priority=800
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cloakbrowser.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:cloakbrowser-cdp-proxy]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",CLOAKBROWSER_CDP_PORT="%(ENV_CLOAKBROWSER_CDP_PORT)s",CLOAKBROWSER_PROXY_PORT="%(ENV_CLOAKBROWSER_PROXY_PORT)s"
command=/usr/local/bin/python3 /usr/local/bin/neko-cloakbrowser-cdp-proxy.py
autorestart=true
priority=805
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cloakbrowser-cdp-proxy.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:cloakbrowser-cdp-init-script]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",CLOAKBROWSER_CDP_PORT="%(ENV_CLOAKBROWSER_CDP_PORT)s",CLOAKBROWSER_DESKTOP_GUARD_ENABLED="%(ENV_CLOAKBROWSER_DESKTOP_GUARD_ENABLED)s",CLOAKBROWSER_DESKTOP_GUARD_SCRIPT="%(ENV_CLOAKBROWSER_DESKTOP_GUARD_SCRIPT)s",CLOAKBROWSER_CDP_INIT_INTERVAL="%(ENV_CLOAKBROWSER_CDP_INIT_INTERVAL)s",CLOAKBROWSER_CDP_INIT_STARTUP_DELAY="%(ENV_CLOAKBROWSER_CDP_INIT_STARTUP_DELAY)s"
command=/usr/local/bin/python3 /usr/local/bin/neko-cloakbrowser-cdp-init-script.py
autorestart=true
priority=806
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cloakbrowser-cdp-init-script.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:cliproxyapi]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",CLIPROXY_HOST="%(ENV_CLIPROXY_HOST)s",CLIPROXY_PORT="%(ENV_CLIPROXY_PORT)s",CLIPROXY_UPSTREAM_BASE_URL="%(ENV_OPENAI_BASE_URL)s",CLIPROXY_UPSTREAM_API_KEY="%(ENV_OPENAI_API_KEY)s",CLIPROXY_API_KEY="%(ENV_CLIPROXY_API_KEY)s",CLIPROXY_PROXY_SERVER="%(ENV_CLOAKBROWSER_PROXY_SERVER)s",CLIPROXY_PROXY_USERNAME="%(ENV_CLOAKBROWSER_PROXY_USERNAME)s",CLIPROXY_PROXY_PASSWORD="%(ENV_CLOAKBROWSER_PROXY_PASSWORD)s",CLIPROXY_IPINFO_URL="%(ENV_CLIPROXY_IPINFO_URL)s"
command=/usr/local/bin/python3 /usr/local/bin/neko-cliproxyapi
autorestart=true
priority=811
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cliproxyapi.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:openbox]
environment=TZ="%(ENV_TZ)s",HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s"
command=/usr/bin/openbox --config-file /etc/neko/openbox.xml
autorestart=true
priority=300
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/openbox.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true
